# userdir.timer
[Unit]
Description=timer for Calculate the size of home directories

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
