# userdir.timer
[Unit]
Description=timer for Calculate the size of home directories

[Timer]
OnCalendar=*-*-* */3:00:00

[Install]
WantedBy=timers.target
